# lmslog.service
[Unit]
Description=Keep LMS up and do logs

[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c "source /storage/.kodi/docker/lmsup.sh"